# resources/feh.service
[Unit]
Description=change bg by feh every 30s

[Service]
Type=simple

